# 90-dd-one-from-udev.rules
KERNEL=="sr[0-9]*", ENV{DEVTYPE}=="disk", ENV{ID_CDROM_MEDIA_DVD}=="1", RUN+="/home/mgalgs/scripts/dd-one-from-udev.sh"
